# systemd/codex-telegram-bot.service.template
[Unit]
Description=Codex Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__PYTHON_BIN__ __HOOK_PY__ bot
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
